# roles/signalk-npm/templates/signalk.socket.j2
[Socket]
ListenStream={{ signalk_server_port }}

[Install]
WantedBy=sockets.target
